# packages/onedrive/root/usr/local/stow/onedrive.service/lib/systemd/user/onedrive.service
[Unit]
Description=OneDrive Free Client
Documentation=https://github.com/skilion/onedrive

[Service]
ExecStart=/usr/bin/trickle -s -u 1000 -d 1500 /usr/bin/onedrive -m
Restart=no

[Install]
WantedBy=default.target
